# ansible/roles/app/templates/aitesters-reset.timer.j2
[Unit]
Description=Daily reset for aitesters sandbox

[Timer]
OnCalendar={{ aitesters_reset_schedule }}
Persistent=true
Unit=aitesters-reset.service

[Install]
WantedBy=timers.target
